# deploy/quadlet/trellis-nats.container
[Unit]
Description=Trellis local NATS
After=network-online.target
Wants=network-online.target

[Container]
ContainerName=trellis-nats
Image=docker.io/library/nats:latest
Exec=-c /etc/nats/nats.conf
PublishPort=127.0.0.1:4222:4222
PublishPort=127.0.0.1:8222:8222
Volume=%h/trellis/config/nats.conf:/etc/nats/nats.conf:ro,Z
Volume=%h/trellis/config/jwt.conf:/etc/nats/jwt.conf:ro,Z
Volume=trellis-nats-data.volume:/data
Volume=trellis-nats-jwt.volume:/jwt
Network=trellis.network
UserNS=keep-id
NoNewPrivileges=true
DropCapability=all
ReadOnly=true
Tmpfs=/tmp:rw,noexec,nosuid,nodev,size=16m

[Service]
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
