# deploy/quadlet/trellis.container
[Unit]
Description=Trellis runtime
After=network-online.target
Wants=network-online.target

[Container]
ContainerName=trellis
Image=ghcr.io/REPLACE_ME/trellis-runtime:REPLACE_ME
PublishPort=3000:3000
EnvironmentFile=/etc/trellis/trellis.env
Volume=/var/lib/trellis/creds:/run/trellis/creds:ro,z
ReadOnly=true
Tmpfs=/tmp

[Service]
Restart=always

[Install]
WantedBy=default.target
